# dynhost_ovh_livebox.timer
[Unit]
Description=Run OVH DynHost update

[Timer]
OnBootSec=1min
OnUnitActiveSec=1days

[Install]
WantedBy=timers.target