# dynhost_ovh_livebox.service
[Unit]
Description=Run OVH DynHost update

[Service]
Type=simple
ExecStart=/opt/dynhost_ovh_livebox/dynhost_ovh_livebox.sh
